# configs/99-digirig-audio.rules
# ReticulumHF: Auto-configure CM108/Digirig audio when USB device is plugged in
#
# This rule triggers when a C-Media CM108 USB audio device (used by Digirig)
# is detected. It runs a script to set optimal audio levels for digital modes.
#
# Vendor ID 0d8c = C-Media Electronics
# Product IDs: 013c (CM108), 0012, 000c, etc.

# When CM108 USB audio device is added, run the audio setup script
ACTION=="add", SUBSYSTEM=="sound", ATTRS{idVendor}=="0d8c", RUN+="/opt/reticulumhf/scripts/set-digirig-audio.sh %n"

# Also match "USB PnP Sound Device" description
ACTION=="add", SUBSYSTEM=="sound", ATTR{id}=="Device", RUN+="/opt/reticulumhf/scripts/set-digirig-audio.sh %n"
